# meta-mitaccomputing/common/core/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-host-bios@.service
[Unit]
Description=Flash Host Bios image %I to Host (MiTAC)

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/bin/sh /usr/sbin/bios-update.sh /tmp/images/%i
